# trading-backend/mt5-bridge/supervisord.conf
; ════════════════════════════════════════════════════════════════════════════
; Supervisor config for mt5-bridge
; - xvfb : virtual display so Wine/MT5 can run "headless"
; - x11vnc : VNC server on :5900 for manual first-time MT5 login
; - mt5 : MT5 terminal process (auto-login via /config.ini)
; ════════════════════════════════════════════════════════════════════════════

[supervisord]
nodaemon=false
logfile=/var/log/supervisord.log
loglevel=info

[program:xvfb]
command=/usr/bin/Xvfb :99 -screen 0 1280x1024x24 -ac +extension GLX +render -noreset
autostart=true
autorestart=true
priority=1
stdout_logfile=/var/log/xvfb.log
stderr_logfile=/var/log/xvfb.err

[program:x11vnc]
command=/usr/bin/x11vnc -display :99 -nopw -listen 0.0.0.0 -rfbport 5900 -forever -shared -quiet
autostart=true
autorestart=true
priority=2
stdout_logfile=/var/log/x11vnc.log
stderr_logfile=/var/log/x11vnc.err
depends_on=xvfb

[program:mt5]
; Wait 5s for Xvfb to settle, then launch MT5 terminal with auto-login config
; /config.ini is generated by entrypoint.sh from MT5_LOGIN/MT5_PASSWORD/MT5_SERVER
; Fusion-branded install creates folder "Fusion Markets MetaTrader 5" — fallback to generic path if absent
; config.ini is placed INSIDE the MT5 folder by entrypoint.sh.
; We pass the bare filename "config.ini" — MT5 resolves it relative to terminal64.exe's
; directory (portable mode), which avoids the Wine Unix→Windows path translation problem.
command=/bin/bash -c "sleep 5 && if [ -f '/root/.wine/drive_c/Program Files/Fusion Markets MetaTrader 5/terminal64.exe' ]; then cd '/root/.wine/drive_c/Program Files/Fusion Markets MetaTrader 5' && wine terminal64.exe /config:config.ini /portable; else cd '/root/.wine/drive_c/Program Files/MetaTrader 5' && wine terminal64.exe /config:config.ini /portable; fi"
autostart=true
autorestart=true
priority=10
stdout_logfile=/var/log/mt5.log
stderr_logfile=/var/log/mt5.err
environment=DISPLAY=":99",WINEPREFIX="/root/.wine",WINEDEBUG="-all"
stopsignal=TERM
stopwaitsecs=10
